# systemd/pqc-sig-charts-restart.path
[Unit]
Description=Watch /home/pqc/pqc-sig-charts-website and restart service on changes

[Path]
PathModified=/home/pqc/pqc-sig-charts-website

[Install]
WantedBy=multi-user.target