# systemd/pqc-sig-charts.service
[Unit]
Description=PQC Digital Signatures website visualization

[Service]
WorkingDirectory=/home/pqc/pqc-sig-charts-website
ExecStart=/home/pqc/pqc-sig-charts-website/run.sh
Restart=always
RestartSec=10
User=pqc
Group=pqc

[Install]
WantedBy=multi-user.target
