# systemd/pqc-sig-charts-debounce.path
[Unit]
Description=Watch PQC Digital Singature webapp dir and debounce restarts

[Path]
PathModified=/home/pqc/pqc-sig-charts-website/.deploy-stamp
Unit=pqc-sig-charts-debounce.service

[Install]
WantedBy=multi-user.target